# systemd/user/checkbattery.service
[Unit]
Description=Script for checking the battery level

[Service]
Type=simple
ExecStart=/home/dennis/.config/i3/checkbattery.sh
